# s6/services.d/02-manager/run
#!/usr/bin/with-contenv bash
source /ctx/bin/functions.sh

logging "Start ICON2 Manager - ${BASE_DIR}, ${GOLOOP_NODE_SOCK}, ${LOGDIR}, ${BOOTING_LOG}"

python3 /ctx/main.py
